# daemon/ducky-game-hub-gamepad.service
[Unit]
Description=Ducky Game Hub Gamepad Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 WORKSPACE_PATH/daemon/ducky_game_hub_gamepad_daemon.py
Restart=on-failure
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
