# daemon/retrocloud-gamepad.service
[Unit]
Description=RetroCloud Gamepad Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 WORKSPACE_PATH/daemon/retrocloud_gamepad_daemon.py
Restart=on-failure
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
